# root/etc/service/cron/run
#!/bin/sh

exec crond -fL /var/log/cron.log
